# src/php/fpm/etc/services.d/10-php-fpm/run
#!/command/with-contenv bash

# Switch to web user
su webuser

#Bring up PHP FPM
/usr/sbin/php-fpm${PHP_VERSION} --nodaemonize