# src/php/fpm/etc/services.d/10-php-fpm/run
#!/command/with-contenv bash

echo
/usr/sbin/php-fpm -v

su webuser
exec /usr/sbin/php-fpm --nodaemonize