# src/php/fpm-nginx/etc/services.d/80-nginx/run
#!/command/with-contenv bash

echo
/usr/sbin/nginx -v
/usr/sbin/nginx -t

echo
exec /usr/sbin/nginx