# src/webapps/crater/root/etc/cron.d/crater-schedule
MAILTO=""
SHELL=/bin/bash
* * * * * webuser /usr/bin/artisan schedule:run $CRATER_SCHEDULE_OPTIONS >/proc/1/fd/1 2>/proc/1/fd/2
