# src/php/fpm-apache/etc/services.d/80-apache2/run
#!/command/with-contenv bash

echo
/usr/sbin/apache2ctl -v
/usr/sbin/apache2ctl configtest

echo
exec /usr/sbin/apache2ctl -DFOREGROUND