# src/webapps/statamic/root/etc/cron.d/statamic-schedule
MAILTO=""
SHELL=/bin/bash
* * * * * webuser /usr/bin/artisan schedule:run $STATAMIC_SCHEDULE_OPTIONS >/proc/1/fd/1 2>/proc/1/fd/2
