# src/servers/fpm-apache/root/etc/services.d/200-apache2/run
#!/command/with-contenv bash
timeout 10 bash -c -- '\
while [ ! -e "${PHPFPM_SOCK_PATH}" ]; do sleep 0.1; done;\
echo "Found ${PHPFPM_SOCK_PATH}!"'
exec apache2ctl -DFOREGROUND