# src/servers/fpm-apache/root/etc/services.d/200-apache2/run
#!/command/with-contenv bash
wait-for-fpm 10 0.1 && exec apache2ctl -DFOREGROUND