# src/servers/fpm-apache/root/etc/services.d/200-apache2/run
#!/command/with-contenv bash
exec apache2ctl -DFOREGROUND