# src/php/fpm-apache/etc/services.d/apache2/run
#!/command/with-contenv sh

exec /usr/sbin/apache2ctl -DFOREGROUND