# src/php/fpm-apache/etc/services.d/20-apache2/run
#!/command/with-contenv sh
/usr/sbin/apache2 -v
exec /usr/sbin/apache2ctl -DFOREGROUND