# src/webapps/hyperf/root/etc/services.d/300-hyperf/run
#!/command/with-contenv bash
wait-for-composer && wait-for-server
echo; echo "🤖 Starting Hyperf Server"
exec wdo hyperf start