# src/php/fpm-nginx/etc/services.d/20-nginx/run
#!/command/with-contenv sh
echo && /usr/sbin/nginx -v
exec /usr/sbin/nginx