# src/php/fpm/root/etc/services.d/100-php-fpm/run
#!/command/with-contenv bash

echo
/usr/sbin/php-fpm -v

echo
su webuser
/usr/sbin/php-fpm --nodaemonize